# imageright-proxy/imageright-proxy.service
[Unit]
Description=ImageRight API proxy for Supabase Edge Functions
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=imageright
Group=imageright
WorkingDirectory=/opt/imageright-proxy
EnvironmentFile=/etc/imageright-proxy/.env
# SOAP transport is plain HTTP (port 80) inside the customer VPN — no cert.
# The legacy NODE_TLS_REJECT_UNAUTHORIZED=0 from the REST era is no longer
# needed and has been removed.
#
# Default Node heap (~2 GB) is fine. The PDF-fetch path streams base64 → decode
# → response chunk-by-chunk via streamPdfForPages (soap.js), so memory is flat
# regardless of PDF size. If you find a `--max-old-space-size` flag here in a
# future hot-patch, the streaming path probably regressed — investigate before
# bumping the heap.
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/imageright-proxy

[Install]
WantedBy=multi-user.target
